# root/etc/s6-overlay/s6-rc.d/acw-auto-library/up
/etc/s6-overlay/s6-rc.d/acw-auto-library/run